# server/ccs.service
[Unit]
Description=Case LED/LCD control service

[Service]
Type=simple
ExecStart=/usr/local/bin/ccs /etc/ccs

[Install]
WantedBy=multi-user.target
